# debian/debian/xnas.target
[Unit]
Description=X-NAS Target
Requires=multi-user.target
After=multi-user.target
AllowIsolate=yes
